# gunicorn/gunicorn.conf
description "Gunicorn application server handling myproject"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid deploy
setgid www-data
chdir /home/deploy

exec .enviroment/bin/gunicorn --workers 3 --bind unix:/home/deploy/weedmatch.sock weedmatch.wsgi:application